# libs/rm2fb/rm2fb.socket
[Unit]
Description=rm2 framebuffer server sockets

[Socket]
ListenDatagram=/var/run/rm2fb.sock
ListenStream=8888

[Install]
WantedBy=sockets.target
