# libs/rm2fb/opkg/postinst
#!/bin/sh
set -eu

# Enable the rm2fb socket
systemctl daemon-reload
systemctl enable --now rm2fb.socket
