# libs/rm2fb/opkg/prerm
#!/bin/sh
set -eu

# Disable the rm2fb socket
systemctl disable --now rm2fb.socket rm2fb.service
